# packaging/pit-panel.service
[Unit]
Description=pit-panel VPS Management Panel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pit-panel
Group=pit-panel
WorkingDirectory=/opt/pit-panel
Environment="PATH=/opt/pit-panel/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PITPANEL_CONFIG_PATH=/etc/pit-panel/config.toml"
ExecStartPre=/bin/mkdir -p /var/lib/pit-panel /opt/pit-panel/apps /var/log/pit-panel
ExecStartPre=/bin/chown pit-panel:pit-panel /var/lib/pit-panel /opt/pit-panel/apps /var/log/pit-panel
ExecStartPre=-/usr/bin/sudo -n /usr/sbin/usermod -a -G systemd-journal pit-panel
ExecStartPre=-/usr/bin/sudo -n /usr/sbin/usermod -a -G docker pit-panel
ExecStart=/opt/pit-panel/.venv/bin/python -c "from pit_panel.main import main; main()"
Restart=always
RestartSec=10

# Sandboxing (compatible with Python + sudo for upgrades)
# Note: NoNewPrivileges is NOT set because pit-panel uses sudo -n
# for the upgrade flow (git fetch/reset, cp service file, systemctl restart).
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictRealtime=true
LockPersonality=true

# Filesystem access
ReadWritePaths=/opt/pit-panel/.venv /opt/pit-panel /etc/systemd/system /var/lib/pit-panel /opt/pit-panel/apps /var/log/pit-panel /etc/pit-panel /etc/caddy /var/lib/apt /var/cache/apt /var/lib/dpkg

# Capabilities: kept empty intentionally so sudo can operate.
# CapabilityBoundingSet= empty would block sudo from changing UID/GID.
# AmbientCapabilities= empty is harmless but unnecessary with sudo.

[Install]
WantedBy=multi-user.target
